# skills/host-your-agent/templates/com.youragent.overnight.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.youragent.overnight</string>

  <!-- caffeinate keeps a sleeping Mac awake for the duration of the run
       (up to the -t seconds below) so a scheduled overnight job actually fires. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-dimsu</string>
    <string>-t</string>
    <string>__CAFFEINATE_SECONDS__</string>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>__RUNNER__</string>
    <string>--job</string>
    <string>__JOB__</string>
    <string>--minutes</string>
    <string>__MINUTES__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin:__HOME__/.npm-global/bin:__HOME__/.bun/bin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>__HOUR__</integer>
    <key>Minute</key>
    <integer>__MINUTE__</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/youragent-overnight/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/youragent-overnight/stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
